# setup/roles/rov/files/rc.local
#! /bin/bash 

ethtool --set-eee eth0 eee off
